# systemd.service
# Systemd service file untuk Pterodactyl Telegram Bot
# Copy file ini ke /etc/systemd/system/pterodactyl-bot.service

[Unit]
Description=Pterodactyl Telegram Control Bot
After=network.target mysql.service
Wants=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/path/to/pterodactyl-telegram-bot
ExecStart=/usr/bin/php /path/to/pterodactyl-telegram-bot/index.php?mode=polling
Restart=always
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

# Environment file
EnvironmentFile=/path/to/pterodactyl-telegram-bot/.env

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/path/to/pterodactyl-telegram-bot/logs

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pterodactyl-bot

[Install]
WantedBy=multi-user.target

# Commands untuk manage service:
# sudo systemctl daemon-reload
# sudo systemctl enable pterodactyl-bot
# sudo systemctl start pterodactyl-bot
# sudo systemctl stop pterodactyl-bot
# sudo systemctl restart pterodactyl-bot
# sudo systemctl status pterodactyl-bot
# sudo journalctl -u pterodactyl-bot -f
